# deploy/roles/scxmppd/templates/scxmppd.service.j2
[Unit]
Description=scxmppd service

[Service]
Restart=always
ExecStart=/bin/bash -c "export ENV={{ env }} && /usr/bin/scxmppd"

[Install]
WantedBy=multi-user.target
